# tools/rebalance_to_gpu.service
[Unit]
Description=Rebalance pods toward racoon-gpu on boot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=KUBECONFIG=/etc/kubernetes/admin.conf
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env bash /home/dante/racoon_server/tools/rebalance_to_gpu.sh
User=root
Group=root
TimeoutStartSec=5min

[Install]
WantedBy=multi-user.target