# tools/rebalance_to_gpu.service
[Unit]
Description=Rebalance Kubernetes pods to racoon-gpu on boot
After=network-online.target kubelet.service crio.service
Wants=network-online.target kubelet.service

[Service]
Type=oneshot
Environment=KUBECONFIG=/etc/kubernetes/admin.conf
ExecStart=/usr/bin/env bash /home/dante/racoon_server/tools/rebalance_to_gpu.sh
User=root
Group=root
RemainAfterExit=true

[Install]
WantedBy=multi-user.target


